# service-name-backup.service
[Unit]

[Service]
Type=oneshot
ExecStart=/opt/backup-db.sh
